# infra/agency/systemd/agency-analytics-backup.timer
[Unit]
Description=Run the Agency Analytics backup nightly

[Timer]
OnCalendar=*-*-* 03:15:00 UTC
RandomizedDelaySec=30m
Persistent=true
Unit=agency-analytics-backup.service

[Install]
WantedBy=timers.target

