# infra/agency/systemd/agency-analytics-backup.service
[Unit]
Description=Agency Analytics encrypted database backup
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/srv/agency-analytics
ExecStart=/srv/agency-analytics/infra/agency/backup.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

